# crypt.service
[Unit]
Description=Decryption

[Service]
ExecStart=/home/henna/bin/crypt.sh

[Install]
WantedBy=multi-user.target
