# decrypt.service
[Unit]
Description=Decryption

[Service]
ExecStart=/home/henna/bin/decrypt.sh

[Install]
WantedBy=multi-user.target
